# dot_config/profile.d/20-util-env.nnn.sh
if command -v nnn >/dev/null 2>&1; then
  export NNN_OPTS="acEoru"
  export NNN_PLUG='p:preview-tabbed;d:dragdrop'
fi
