# dot_config/profile.d/20-util-env.terminfo.sh
if [ -d "$HOME/.config/terminfo" ]; then
  export TERMINFO="$HOME/.config/terminfo"
fi
