# dot_config/profile.d/20-util-env.lftp.sh
if [ -d "$HOME/.config/lftp" ]; then
  export LFTP_HOME="$HOME/.config/lftp"
fi
