# dot_config/profile.d/20-util-env.exa.sh
# via https://github.com/ogham/exa/issues/280#issuecomment-959824115
export EXA_COLORS=di=34:bd=33:cd=33:so=31:ex=32:ur=33:uw=31:ux=32:ue=32:uu=33:gu=33:lc=31:df=32:sn=32:nb=32:nk=32:nm=32:ng=32:nt=32
