# dot_local/share/systemd/user/ffs-flush.service
[Unit]
Description=Firefox session flush service

[Service]
ExecStart=/home/rlue/.local/bin/ffs flushall
Restart=no
